# media-tv/oscam/files/oscam
#!/sbin/runscript
# Copyright 1999-2009 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: Exp $

depend() {
	before vdr
}

start() {
	ebegin "Starting OSCam cardserver"
	/usr/bin/oscam -b -c /etc/oscam
	eend $?
}

stop() {
	ebegin "Stopping OSCam cardserver"
	start-stop-daemon --stop --quiet --pidfile /var/run/oscam.pid
	eend $?
}
